# Diskspace/diskspacemon.timer
[Unit]
Description=Run diskspacemon once a day

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
